# debian/postinst
#! /bin/sh
# postinst script for tinc
#
# $Id: postinst,v 1.6 2000/05/21 23:01:28 zarq Exp $
#
# see: dh_installdeb(1)

TCONF="/etc/tinc"
NETSFILE="$TCONF/nets.boot"

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see /usr/doc/packaging-manual/
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.

case "$1" in
    configure)
	if [ ! -e /dev/.devfsd -a ! -e /dev/.devfs ] ; then
	    devices_exist=1
	    for d in 0 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15 ; do
		test -e /dev/tap$d || devices_exist=0
	    done
	    if [ $devices_exist -eq 0 ] ; then
	        cat <<EOF

You don't seem to have the necessary devices available.  If you want,
I can create them for you now, or you may create them later yourself
with the command:

  cd /dev && ./MAKEDEV netlink

EOF
		echo -n 'Do you want me to do this for you now? '
		read answer
		case $answer in
		    y*|Y*)
			cd /dev && ./MAKEDEV netlink || true
			;;
		    *)
			echo 'Ok.  But you must create these devices before tinc can work!'
			;;
		esac
	    fi
	fi

	if [ ! -e $NETSFILE ] ; then
	    echo "## This file contains all names of the networks to be started on system startup." > $NETSFILE
	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# Automatically added by dh_installinit, edited for use with debconf
if [ -x "/etc/init.d/tinc" ]; then
    update-rc.d tinc defaults >/dev/null

    if [ "$1" = "configure" ] ; then
	if dpkg --compare-versions "$2" lt "1.0pre5-4" ; then
	    /etc/init.d/tinc start
	else
	    if [ -e /usr/share/debconf/confmodule ] ; then
		. /usr/share/debconf/confmodule
		db_get tinc/restart_on_upgrade
		if [ "$RET" = "true" ] ; then
		    /etc/init.d/tinc start
		else
		    echo "I have not stopped or restarted the tinc daemon.  You should do this"
		    echo "yourself whenever you're ready; type \`\`/etc/init.d/tinc restart''."
		fi
		db_stop
	    else
		/etc/init.d/tinc start
	    fi
	fi
    fi
fi
# End automatically added section

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


